# scripts/getty@tty1.service
[Service]
ExecStart=
ExecStart=-/sbin/agetty -o '-p -f -- \\u' --noclear --autologin ukko %I $TERM
Type=simple
Environment=XDG_SESSION_TYPE=x11
